# tools/archive-proxy/kgmexporter-proxy.service
[Unit]
Description=kgmexporter archive.org upload proxy
After=network.target

[Service]
Type=simple
User=archive
WorkingDirectory=/opt/kgmexporter-proxy
# Real IA_ACCESS / IA_SECRET / IA_ITEM / PROXY_TOKEN live in /etc/kgmexporter-proxy.env
# (chmod 600, owner root) so they are never checked into the repo.
EnvironmentFile=/etc/kgmexporter-proxy.env
ExecStart=/opt/kgmexporter-proxy/venv/bin/gunicorn -w 4 --worker-class gthread --threads 2 -b 0.0.0.0:9080 --timeout 1800 proxy:app
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
